# recipes/files/jvideo-queue-monitor-python.service
[Unit]
Description=Juni's Queue Monitor Service (Python)
After=network.target redis.service
Wants=redis.service
Conflicts=jvideo-queue-monitor-cpp.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jvideo
Environment="PYTHONPATH=/opt/jvideo/services"
ExecStart=/usr/bin/python3 /opt/jvideo/services/queue_monitor.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
