# recipes/files/jvideo-frame-saver-python.service
[Unit]
Description=Juni's Frame Saver Service (Python)
After=network.target
Conflicts=jvideo-frame-saver-cpp.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jvideo
Environment="PYTHONPATH=/opt/jvideo/services"
Environment="JVIDEO_SERVICE_IMPL=python"
ExecStart=/usr/bin/python3 /opt/jvideo/services/frame_saver.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
